# crontab
00 18 * * * g213-led -a 00ff00
00 10 * * * g213-led -a 000000 
